# bringup/system/root-init2.service
[Unit]
Description=Root Init 2 Script Runner
After=home.mount

[Service]
ExecStart=/home/root/init2.sh

[Install]
WantedBy=multi-user.target
